# var/lib/dpkg/info/ipmitool.conffiles
/etc/init.d/ipmievd
/etc/default/ipmievd
