# var/lib/dpkg/info/ipmitool.list
/.
/etc
/etc/init.d
/etc/init.d/ipmievd
/etc/default
/etc/default/ipmievd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ipmievd.service
/usr
/usr/sbin
/usr/sbin/ipmievd
/usr/bin
/usr/bin/ipmitool
/usr/share
/usr/share/doc
/usr/share/doc/ipmitool
/usr/share/doc/ipmitool/README.gz
/usr/share/doc/ipmitool/AUTHORS
/usr/share/doc/ipmitool/copyright
/usr/share/doc/ipmitool/changelog.Debian.gz
/usr/share/doc/ipmitool/changelog.gz
/usr/share/ipmitool
/usr/share/ipmitool/oem_ibm_sel_map
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ipmievd.8.gz
/usr/share/man/man1
/usr/share/man/man1/ipmitool.1.gz
/usr/lib
/usr/lib/modules-load.d
/usr/lib/modules-load.d/ipmievd.conf
