# var/lib/dpkg/info/nginx.list
/.
/usr
/usr/sbin
/usr/sbin/nginx
/usr/share
/usr/share/doc
/usr/share/doc/nginx
/usr/share/doc/nginx/changelog.gz
/usr/share/doc/nginx/CHANGES.ru.gz
/usr/share/doc/nginx/copyright
/usr/share/doc/nginx/changelog.Debian.gz
/usr/share/doc/nginx/README
/usr/share/nginx
/usr/share/nginx/html
/usr/share/nginx/html/50x.html
/usr/share/nginx/html/index.html
/var
/var/log
/var/log/nginx
/var/cache
/var/cache/nginx
/etc
/etc/logrotate.d
/etc/logrotate.d/nginx
/etc/default
/etc/default/nginx
/etc/nginx
/etc/nginx/koi-win
/etc/nginx/conf.d
/etc/nginx/conf.d/default.conf
/etc/nginx/conf.d/example_ssl.conf
/etc/nginx/nginx.conf
/etc/nginx/mime.types
/etc/nginx/scgi_params
/etc/nginx/fastcgi_params
/etc/nginx/uwsgi_params
/etc/nginx/koi-utf
/etc/nginx/win-utf
/etc/init.d
/etc/init.d/nginx
