# var/lib/dpkg/info/opendkim.conffiles
/etc/default/opendkim
/etc/init.d/opendkim
/etc/opendkim.conf
