# var/lib/dpkg/info/opendmarc.list
/.
/etc
/etc/init.d
/etc/init.d/opendmarc
/etc/default
/etc/default/opendmarc
/etc/opendmarc.conf
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/opendmarc.8.gz
/usr/share/man/man8/opendmarc-reports.8.gz
/usr/share/man/man8/opendmarc-params.8.gz
/usr/share/man/man8/opendmarc-importstats.8.gz
/usr/share/man/man8/opendmarc-import.8.gz
/usr/share/man/man8/opendmarc-expire.8.gz
/usr/share/man/man8/opendmarc-check.8.gz
/usr/share/man/man5
/usr/share/man/man5/opendmarc.conf.5.gz
/usr/share/doc
/usr/share/doc/opendmarc
/usr/share/doc/opendmarc/README.gz
/usr/share/doc/opendmarc/opendmarc.conf.sample.gz
/usr/share/doc/opendmarc/changelog.Debian.gz
/usr/share/doc/opendmarc/changelog.gz
/usr/share/doc/opendmarc/README.opendmarc
/usr/share/doc/opendmarc/copyright
/usr/share/doc/opendmarc/README.Debian
/usr/share/doc/opendmarc/schema.mysql
/usr/share/doc/opendmarc/README.schema
/usr/share/doc/opendmarc/mkdb.mysql
/usr/sbin
/usr/sbin/opendmarc
/usr/sbin/opendmarc-check
/usr/sbin/opendmarc-reports
/usr/sbin/opendmarc-params
/usr/sbin/opendmarc-importstats
/usr/sbin/opendmarc-import
/usr/sbin/opendmarc-expire
